# deploy/systemd/trading-bot.service
[Unit]
Description=Solana Revolutionary Trading Bot
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=trader
Group=trader
WorkingDirectory={{REPO_ROOT}}
EnvironmentFile={{REPO_ROOT}}/.env
ExecStart=/usr/bin/python3 scripts/run_bot.py --network mainnet --no-auto-trade
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
